# skolo/raspberry-pi/rack4/skolo-setup.service
[Unit]
Description=Minimal Skolo Configuration

[Service]
Type=oneshot
ExecStart=/opt/skolo/setup.sh

[Install]
RequiredBy=basic.target
