# skolo/ansible/switchboard/test_machine_monitor.service
[Unit]
Description=test_machine_monitor
After=syslog.target network.target

[Service]
Type=simple
User=chrome-bot
ExecStart=/usr/local/bin/test_machine_monitor \
   --config=prod.json \
   --username=chrome-bot
Restart=always

[Install]
WantedBy=multi-user.target