# skolo/sys/metadata-server.service
[Unit]
Description=Serve metadata to machines in the lab.
Requires=network-online.target
After=network-online.target

[Service]
ExecStart=/usr/local/bin/metadata-server \
  --logtostderr \
  --port=:8001 \
  --prom_port=:20003 \
  --token_file=/var/local/token_swarming.json

Restart=always
User=chrome-bot
Group=chrome-bot

[Install]
WantedBy=multi-user.target
