# datahopper/sys/buildbot-db-backup.timer
[Unit]
Description=Run buildbot_db_backup.service nightly.

[Timer]
OnCalendar=*-*-* 05:00:00 UTC
Persistent=false

[Install]
WantedBy=timers.target
