# skfe/sys/rotate_nginx_logs.timer
[Unit]
Description=Run log rotation once a week

[Timer]
OnCalendar=weekly
Persistent=true
