# skolo/ansible/switchboard/roles/install_prometheus_and_alert_to_pubsub/templates/alert_to_pubsub.service
[Service]
User=chrome-bot
Group=chrome-bot
ExecStart=/usr/local/bin/alert_to_pubsub_ansible \
  --location={{ inventory_hostname }} \
  --port=:{{ all.prometheus.alert_to_pubsub_port }} \
  --prom_port=:{{ all.prometheus.monitoring.ports.alert_to_pubsub_ansible }}

[Install]
WantedBy=multi-user.target